# crontab
@reboot while true ; do "${HOME}/src/temperature.py" ; sleep 5 ; done
