# crontab
@reboot while true ; do "${HOME}/dht11_mqtt/temperature.py" ; sleep 5 ; done
